# templates/default/wrapper_script.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides:             docker_<%= @service_name %>
# Required-Start:       $remote_fs docker
# Required-Stop:        $remote_fs docker
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    Docker container for <%= @service_name %>.
### END INIT INFO

<% @actions.each_pair do |action, command| %>
<%= action %>_run(){
  <%= command %>
}

<% end %>
case "$1" in
  <% @actions.each_pair do |action, command| %>
  <%= action %>)
    <%= action %>_run
    ;;

  <% end %>
  *)
    echo "Usage: ${0##*/}: { <%= @actions.keys.join(' | ') %> }" 2>&1
    exit 65
    ;;
esac
